# opencpu-server/systemd/cleanocpu.timer
[Unit]
Description=Run cleanocpu.service every 10 minutes

[Timer]
OnCalendar=*:0/10
