# opencpu-server/cron.d/opencpu-bioc
0 2 * * 0 www-data Rscript /usr/lib/opencpu/scripts/updatebioc.R
