# opencpu-server/cron.d/opencpu
*/10 * * * * root /usr/lib/opencpu/scripts/cleanocpu.sh
