# opencpu-server/cron.d/opencpu-cran
0 2 * * * www-data Rscript /usr/lib/opencpu/scripts/updatecran.R
