# opencpu-server/cron.d/opencpu
*/10 * * * * www-data /usr/lib/opencpu/scripts/cleanocpu.sh
