# opencpu-server/cron.d/opencpu-tmp
*/10 * * * * opencpu /usr/lib/opencpu/scripts/cleantmp.sh
