# opencpu-cache/systemd/opencpu-cache.service
[Unit]
Description=OpenCPU Caching Server

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/etc/init.d/opencpu-cache start
ExecStartPost=/bin/systemctl reload nginx
ExecStop=/etc/init.d/opencpu-cache stop
ExecStopPost=/bin/systemctl reload nginx

[Install]
After=nginx.service
